# lib/itamae/plugin/recipe/droonga/templates/lib/systemd/system/droonga-http-server.service.erb
[Unit]
Description=Droonga HTTP Server

[Service]
User=droonga-http-server
ExecStart=/usr/bin/droonga-http-server --cache-size=-1
ExecStop=/bin/kill -TERM $MAINPID
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
WorkingDirectory=/home/droonga-http-server/droonga

[Install]
WantedBy=multi-user.target
